# services/foghorn.service
[Unit]
Description=Foghorn %MODE%
After=network.target

[Service]
ExecStart=%COMMAND%
Restart=on-abort